# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── 1. Virtual display ────────────────────────────────────────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)sx24 -ac +extension GLX +render -noreset
user=agent
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

# ── 2. Chromium with CDP ──────────────────────────────────────────────────────
[program:chromium]
command=/usr/local/bin/start-chrome.sh
user=agent
autostart=true
autorestart=true
priority=200
startretries=3
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=":99",HOME="/home/agent"

# ── 3. VNC server (interactive) ───────────────────────────────────────────────
[program:x11vnc]
command=/usr/local/bin/start-x11vnc.sh
user=agent
autostart=true
autorestart=true
priority=300
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=":99",HOME="/home/agent"

# ── 4. noVNC WebSocket bridge ─────────────────────────────────────────────────
[program:websockify]
command=python3 -m websockify --web /opt/novnc 6080 localhost:5900
user=agent
autostart=true
autorestart=true
priority=400
startretries=5
startsecs=2
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify.log

# ── 5. ttyd (Terminal) ────────────────────────────────────────────────────────
[program:ttyd]
command=ttyd --port 7681 --base-path /terminal --writable bash
user=agent
autostart=true
autorestart=true
priority=450
startsecs=2
stdout_logfile=/var/log/supervisor/ttyd.log
stderr_logfile=/var/log/supervisor/ttyd.log
environment=HOME="/home/agent",USER="agent",TERM="xterm-256color"

# ── 6. code-server (VSCode) ───────────────────────────────────────────────────
[program:code-server]
command=code-server --config /home/agent/.config/code-server/config.yaml --base-path /code %(ENV_WORKSPACE)s
user=agent
autostart=true
autorestart=true
priority=500
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/agent",USER="agent"

# ── 7. nginx reverse proxy ────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
priority=600
startsecs=2
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
